# system_files/usr/lib/systemd/user/arcade-background.timer
[Unit]
Description=Refresh Cosmic Arcade background daily

[Timer]
OnBootSec=2min
OnUnitActiveSec=12h
Persistent=true
Unit=arcade-background.service

[Install]
WantedBy=timers.target
